# src/rpc/config/supervisor/feed_server.conf
[program:feed_server]
command         = /home/yvan/bin/feed_server feed_server_sample.conf
directory       = /home/yvan/bin/
user            = yvan
numprocs        = 1
autostart       = true
autorestart     = true
startsecs       = 3
stdout_logfile  = /var/log/supervisor/feed_server.log
redirect_stderr = true
environment     = LOGPATH=/tmp, LANG=zh_CN.utf8, LC_ALL=zh_CN.UTF-8, LC_LANG=zh_CN.UTF-8,
